# deploy/systemd/contaneta-backup.service
# ContaNeta Nightly Backup
# Triggered by contaneta-backup.timer
# Copy to /etc/systemd/system/contaneta-backup.service

[Unit]
Description=ContaNeta Nightly Backup (DB + storage)
After=network.target

[Service]
Type=oneshot
User=deploy
Group=deploy
WorkingDirectory=/opt/contaneta
EnvironmentFile=/opt/contaneta/.env
ExecStart=/opt/contaneta/scripts/backup_nightly.sh
StandardOutput=append:/opt/contaneta/logs/backup.log
StandardError=append:/opt/contaneta/logs/backup.log

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/contaneta
PrivateTmp=yes
